# fnando/Library/LaunchAgents/com.fnando.backup-brew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fnando.backup-brew</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/fnando/.bin/backup-brew</string>
  </array>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>StandardOutPath</key>
  <string>/tmp/backup-brew.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/backup-brew.log</string>
  <key>KeepAlive</key>
  <dict>
      <key>Crashed</key>
      <true/>
  </dict>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
